# .agent-loop/systemd/agent-loop-scheduler.service
DO NOT EDIT — generated by install.py from systemd/agent-loop-scheduler.service (kit 5f496b2db042). Edit the source and run install.py --update.
[Unit]
Description=Loop engineering next-turn scheduler daemon
After=default.target

[Service]
Type=simple
WorkingDirectory=/home/hito/work/git/loopeng-bootstrap
ExecStart=/usr/bin/env python3 /home/hito/work/git/loopeng-bootstrap/.agent-loop/bin/next_turn_scheduler_daemon.py --repo /home/hito/work/git/loopeng-bootstrap
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
